# contrib/aur/common/bitcoin-reindex.service
[Unit]
Description=Bitcoin Daemon (reindex)
Conflicts=bitcoin.service
After=network.target

[Service]
User=bitcoin
Group=bitcoin
ExecStart=/usr/bin/doged -daemon=0                       \
                            -conf=/etc/bitcoin/dogecoin.conf \
                            -datadir=/srv/bitcoin           \
                            -reindex                        \
                            -pid=/run/bitcoin/doged.pid
ExecReload=/usr/bin/kill -HUP $MAINPID
ExecStop=/usr/bin/doged stop
PrivateTmp=true

[Install]
WantedBy=multi-user.target
